# docker/cms-services/rootfs/etc/s6-overlay/s6-rc.d/cmsLogService/run
#!/command/with-contenv bash

exec cmsLogService
