# docker/cms-services/rootfs/etc/s6-overlay/s6-rc.d/cmsInitDB/up
/etc/s6-overlay/s6-rc.d/cmsInitDB/s6-cmsInitDB