# docker/cms-services/rootfs/etc/s6-overlay/s6-rc.d/cmsInitDB/up
#!/command/with-contenv bash

if $CMS_INIT; then
  exec cmsInitDB
fi

