# docker/cms-services/rootfs/etc/s6-overlay/s6-rc.d/cmsImportContest/up
/etc/s6-overlay/s6-rc.d/cmsImportContest/s6-cmsImportContest