# docker/services/rootfs/etc/s6-overlay/s6-rc.d/cmsResourceService/run
#!/command/execlineb -P
cmsResourceService -a $CONTEST_ID
